# service/log/run
#!/bin/sh
exec 2>&1
exec setuidgid kurain multilog t ./main

